# ansible/roles/webserver/files/consumer.service
[Unit]
Description=Kafka Consumer Service

[Service]
ExecStart=/usr/bin/go run consumer.go
WorkingDirectory=/home/ec2-user/
Restart=always
