# ansible/roles/kafka/files/producer.service
[Unit]
Description=Kafka Producer Service

[Service]
ExecStart=/usr/bin/go run producer.go
WorkingDirectory=/home/ec2-user/
